# pitchfork-csv/pitchfork-csv.cron
MAILTO=colin.242@gmail.com
30 13 * * * ec2-user /home/ec2-user/pitchfork-csv/daily-update.rb > /tmp/daily-update-$(date +'\%m\%d\%y\%H\%M').log 2>&1
